# attitude_rootfs_files/etc/crontabs/root
* * * * * /sbin/pi_ping_test
*/2 * * * * /sbin/pi_alive_check

